# cookbooks/fb_logrotate/templates/default/logrotate.service.erb
[Unit]
Description=Rotate log files
Documentation=man:logrotate(8) man:logrotate.conf(5)

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate /etc/logrotate.conf
Nice=<%= node['fb_logrotate']['systemd_settings']['Nice'] %>
IOSchedulingClass=<%= node['fb_logrotate']['systemd_settings']['IOSchedulingClass'] %>
