# cookbooks/fb_swap/templates/default/manage-swap-file.service.erb
[Unit]
Description=Manage swap file
RequiresMountsFor=<%= node['fb_swap']['filesystem'] %>

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/libexec/manage-swap-file
