# cookbooks/fb_swap/templates/default/manage-swap-device.service.erb
[Unit]
Description=Format size of swap device
DefaultDependencies=no
Requires=-.mount
After=-.mount

[Service]
Type=oneshot
RemainAfterExit=yes
<%
require 'shellwords'
device = FB::FbSwap._device(node)
# if there's an existing device, we should round-trip the existing label and
# uuid. This way if this code stops being used, and we fall back on the entry
# in fstab, we'll still find the device using the same predicate as before.
info = node.filesystem_data['by_device'][device]
label = info['label']
uuid = info['uuid']
# While the command line here is systemd's command, not a regular shell this is
# still very close to correct.
options = label ? "--label #{Shellwords.escape(label)} " : ''
options = "#{options}--uuid #{Shellwords.escape(uuid)} " if uuid
-%>
<% if node['fb_swap']['enabled'] -%>
ExecStart=/usr/sbin/mkswap <%= options %><%= device %> <%= node['fb_swap']['_calculated']['device_size_bytes'] / 1024 %>
<% else -%>
ExecStart=/usr/bin/echo "Not formatting swap when disabled"
<% end -%>
